# create_netns.service
[Unit]
Description=A service unit to manage network namespaces.

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/bin/bash /usr/sbin/manage_netns.sh --debug start

[Install]
After=network-online.target
Wants=network-online.target
WantedBy=multi-user.target
